# src/systemd/wsl2-wayland.socket
[Socket]
ListenStream=%t/wayland-0

[Install]
WantedBy=sockets.target