# src/systemd/wsl2-wayland@.socket
[Socket]
ListenStream=/run/user/%i/wayland-0
Requires=user-runtime-dir@.service
After=user-runtime-dir@.service

[Install]
WantedBy=sockets.target