# src/systemd/wsl2-pulseaudio.socket
[Socket]
ListenStream=%t/pulse/native

[Install]
WantedBy=sockets.target