# src/systemd/wsl2-pulseaudio.socket
[Socket]
ListenStream=%t/pulse/native
Requires=user-runtime-dir@.service
After=user-runtime-dir@.service

[Install]
WantedBy=sockets.target