# src/systemd/wsl2-wayland.socket
[Socket]
ListenStream=%t/wayland-0
Requires=user-runtime-dir@.service
After=user-runtime-dir@.service

[Install]
WantedBy=sockets.target